# my_files/luci-app-wifimgr/root/etc/uci-defaults/99-wifimgr-defaults
#!/bin/sh

# Set country to CZ on all radios
uci set wireless.radio0.country='CZ'
uci set wireless.radio1.country='CZ'
uci set wireless.radio2.country='CZ'

# Set default password and rename MT76_AP_MLD on all wifi interfaces
for iface in $(uci show wireless | grep '=wifi-iface' | cut -d'.' -f2 | cut -d'=' -f1); do
    uci set wireless.${iface}.key='12345678'
    ssid=$(uci -q get wireless.${iface}.ssid)
    [ "$ssid" = 'MT76_AP_MLD' ] && uci set wireless.${iface}.ssid='OpenWrt-MLD'
done

uci commit wireless
